# config/etc/systemd/system/renew-certs.service
[Unit]
Description=Renews the TLS certificates

[Service]
Type=oneshot
ExecStart=/usr/bin/renew-certs
ExecStartPost=/usr/bin/systemctl restart nginx
