# config/etc/systemd/system/renew-certs.timer
[Unit]
Description=Run renew-certs.service every month

[Timer]
OnCalendar=monthly

[Install]
WantedBy=multi-user.target
